# ets2-trucknav.service
[Unit]
Description=ETS2 / ATS → TruckNav-Sim telemetry relay and bridge
After=network.target

[Service]
Type=simple
User=<USER>
Group=<USER>
WorkingDirectory=/home/<USER>/ets2-trucknav-linux
ExecStart=/home/<USER>/ets2-trucknav-linux/start.sh
Restart=on-failure
RestartSec=5
Environment="HOME=/home/<USER>"

[Install]
WantedBy=multi-user.target
