# ets2-trucknav.service
[Unit]
Description=ETS2 / ATS → TruckNav-Sim telemetry relay and bridge
After=network.target

[Service]
Type=simple
User=jason
Group=jason
WorkingDirectory=/home/jason/ets2-trucknav-linux
ExecStart=/home/jason/ets2-trucknav-linux/start.sh
Restart=on-failure
RestartSec=5
Environment="HOME=/home/jason"

[Install]
WantedBy=multi-user.target
